# config/supervisord.d/nginx.ini
[program:nginx]
command=/usr/sbin/nginx -c /home/nrpuser/.local/etc/nginx/nginx.conf -g "daemon off;"
numprocs=1
numprocs_start=0
priority=50
autostart=true
autorestart=unexpected
startsecs=1
startretries=3
exitcodes=0,2
stopsignal=TERM
stopwaitsecs=10
killasgroup=false
user=root
redirect_stderr=false
stdout_logfile=/var/log/supervisor/nginx/nginx.out
stdout_logfile_maxbytes=125MB
stdout_logfile_backups=3
stderr_logfile=/var/log/supervisor/nginx/nginx.err
stderr_logfile_maxbytes=125MB
stderr_logfile_backups=3
