# config/supervisord.d/nrp-services_app.ini
[program:nrp-services_app]
command=/bin/bash -c "source $HBP/user-scripts/nrp_variables && PYTHONPATH=$PYTHONPATH:$HOME/.opt/platform_venv/lib/python2.7/site-packages && sleep 2 && exec uwsgi --ini  $HOME/.local/etc/nginx/uwsgi-nrp.ini"
numprocs=1
numprocs_start=0
priority=500
autostart=true
autorestart=unexpected
startsecs=1
startretries=3
exitcodes=0,2
stopsignal=TERM
stopwaitsecs=10
killasgroup=false
environment=HBP=%(ENV_HBP)s
user=bbpnrsoa
redirect_stderr=false
stdout_logfile=/var/log/supervisor/nrp-services_app/nrp-services_app.out
stdout_logfile_maxbytes=125MB
stdout_logfile_backups=3
stderr_logfile=/var/log/supervisor/nrp-services_app/nrp-services_app.err
stderr_logfile_maxbytes=125MB
stderr_logfile_backups=3
