# config/supervisord.d/ros-simulation-factory_app.ini
[program:ros-simulation-factory_app]
command=/bin/bash -c "source $HBP/user-scripts/nrp_variables && export ROS_IP=%(ENV_ROS_IP)s && export PYTHONPATH=$PYTHONPATH:$HOME/.opt/platform_venv/lib/python2.7/site-packages && exec python /home/bbpnrsoa/nrp/src/ExDBackend/hbp_nrp_cleserver/hbp_nrp_cleserver/server/ROSCLESimulationFactory.py"
numprocs=1
numprocs_start=0
priority=400
autostart=true
autorestart=unexpected
startsecs=1
startretries=3
exitcodes=0,2
stopsignal=TERM
stopwaitsecs=10
killasgroup=false
directory=/home/bbpnrsoa/.opt/platform_venv
environment=HBP=%(ENV_HBP)s
user=bbpnrsoa
redirect_stderr=false
stdout_logfile=/var/log/supervisor/ros-simulation-factory_app/ros-simulation-factory_app.out
stdout_logfile_maxbytes=125MB
stdout_logfile_backups=3
stderr_logfile=/var/log/supervisor/ros-simulation-factory_app/ros-simulation-factory_app.err
stderr_logfile_maxbytes=125MB
stderr_logfile_backups=3
