# config/supervisord.d/rosbridge.ini
[program:rosbridge]
command=/bin/bash -c "source $HBP/user-scripts/nrp_variables && exec /opt/ros/kinetic/bin/rosrun rosbridge_server rosbridge_websocket"
numprocs=1
numprocs_start=0
priority=300
autostart=true
autorestart=unexpected
startsecs=1
startretries=3
exitcodes=0,2
stopsignal=TERM
stopwaitsecs=10
killasgroup=false
environment=HBP=%(ENV_HBP)s
user=bbpnrsoa
redirect_stderr=false
stdout_logfile=/var/log/supervisor/rosbridge/rosbridge.out
stdout_logfile_maxbytes=125MB
stdout_logfile_backups=3
stderr_logfile=/var/log/supervisor/rosbridge/rosbridge.err
stderr_logfile_maxbytes=125MB
stderr_logfile_backups=3
