# config/supervisord.f/proxy.ini
[program:proxy_app]
command=/bin/bash -c "source $HBP/user-scripts/nrp_variables && cd /home/bbpnrsoa/nrp/src/nrpBackendProxy && node app.js"
numprocs=1
numprocs_start=0
priority=999
autostart=true
autorestart=unexpected
startsecs=1
startretries=3
exitcodes=0,2
stopsignal=TERM
stopwaitsecs=10
killasgroup=false
environment=HBP=%(ENV_HBP)s
user=bbpnrsoa
redirect_stderr=false
stdout_logfile=/var/log/supervisor/proxy/proxy.out
stdout_logfile_maxbytes=250MB
stdout_logfile_backups=10
stderr_logfile=/var/log/supervisor/proxy/proxy.err
stderr_logfile_maxbytes=250MB
stderr_logfile_backups=10
